# files/etc/profile.d/bash.sh
if [ -f /etc/bash_history_log ]; then
    . /etc/bash_history_log
fi
